# deploy/antfly.service
[Unit]
Description=Antfly vector database
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/home/ubuntu/antfly/antfly swarm --config /home/ubuntu/antfly/gif-picker/config.yaml
Restart=on-failure
RestartSec=5

# Environment
Environment=HOME=/home/ubuntu

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
